# backend/docker/supervisord.conf
[supervisorctl]
serverurl=unix:///var/run/supervisord.sock

[unix_http_server]
file=/var/run/supervisord.sock
chmod=0700

[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10
loglevel=info
pidfile=/var/run/supervisord.pid

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:queue-worker]
command=php /app/artisan queue:work --sleep=1 --tries=3
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/queue.log
stdout_logfile_maxbytes=10MB
stderr_logfile=/var/log/supervisor/queue-error.log
stderr_logfile_maxbytes=10MB
redirect_stderr=false
stopwaitsecs=3600

[program:kafka-consumer]
command=php /app/artisan kafka:billing-consume
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/kafka.log
stdout_logfile_maxbytes=10MB
stderr_logfile=/var/log/supervisor/kafka-error.log
stderr_logfile_maxbytes=10MB
redirect_stderr=false
stopwaitsecs=3600

[program:payment-consumer]
command=php /app/artisan kafka:payment-consume
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/payment.log
stdout_logfile_maxbytes=10MB
stderr_logfile=/var/log/supervisor/payment-error.log
stderr_logfile_maxbytes=10MB
redirect_stderr=false
stopwaitsecs=3600

[program:notification-consumer]
command=php /app/artisan kafka:notification-consume
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/notification.log
stdout_logfile_maxbytes=10MB
stderr_logfile=/var/log/supervisor/notification-error.log
stderr_logfile_maxbytes=10MB
redirect_stderr=false
stopwaitsecs=3600

[program:service-activated-consumer]
command=php /app/artisan kafka:service-activated-consume
autostart=true
autorestart=true
user=root
stdout_logfile=/var/log/supervisor/service-activated.log
stdout_logfile_maxbytes=10MB
stderr_logfile=/var/log/supervisor/service-activated-error.log
stderr_logfile_maxbytes=10MB
redirect_stderr=false
stopwaitsecs=3600

[program:subscription-cancelled-consumer]
command=php /app/artisan kafka:subscription-cancelled-consume
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/subscription-cancelled.log
stdout_logfile_maxbytes=10MB
stderr_logfile=/var/log/supervisor/subscription-cancelled-error.log
stderr_logfile_maxbytes=10MB
redirect_stderr=false
stopwaitsecs=3600

[program:plan-updated-consumer]
command=php /app/artisan kafka:plan-updated-consume
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/plan-updated.log
stdout_logfile_maxbytes=10MB
stderr_logfile=/var/log/supervisor/plan-updated-error.log
stderr_logfile_maxbytes=10MB
redirect_stderr=false
stopwaitsecs=3600

[program:plan-deactivated-consumer]
command=php /app/artisan kafka:plan-deactivated-consume
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/plan-deactivated.log
stdout_logfile_maxbytes=10MB
stderr_logfile=/var/log/supervisor/plan-deactivated-error.log
stderr_logfile_maxbytes=10MB
redirect_stderr=false
stopwaitsecs=3600

[program:plan-deactivated-consumer]
command=php /app/artisan kafka:plan-deactivated-consume
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/plan-deactivated.log
stdout_logfile_maxbytes=10MB
stderr_logfile=/var/log/supervisor/plan-deactivated-error.log
stderr_logfile_maxbytes=10MB
redirect_stderr=false
stopwaitsecs=3600
